# edge-agent/drift-deploy-agent.service
[Unit]
Description=Drift Deploy edge agent
Documentation=https://github.com/kidproquo/drift
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/drift-deploy/env
ExecStart=/usr/local/bin/drift-deploy-agent.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Run as root because the agent calls `docker compose` and writes
# into /var/lib/drift-deploy. Drop privileges in v1 with a `drift` user
# in the `docker` group.

[Install]
WantedBy=multi-user.target
